# supervisor/conf.d/celery.conf
; Celery — worker + beat (-B) dans le même processus. Remplace l'ancien
; conteneur lespass_celery séparé du docker-compose de production
; (un seul conteneur = moins d'empreinte mémoire).
; / Celery — worker + beat (-B) in one process. Replaces the former
; separate lespass_celery container in the production docker-compose
; (single container = smaller memory footprint).
[program:celery]
command=/home/tibillet/.local/bin/poetry run celery -A TiBillet worker -l INFO -B --concurrency=6
directory=/DjangoFiles
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/DjangoFiles/logs/celery.logs
