# supervisor/conf.d/daphne.conf
; Daphne — serveur ASGI dédié aux WebSockets (port 7999) : ws/laboutik/,
; ws/printer/, ws/rfid/ (controlvanne), ws/chat/. nginx route les URLs
; /ws/ et /wss/ vers ce port (voir nginx_prod/lespass_prod.conf).
; / Daphne — ASGI server dedicated to WebSockets (port 7999): ws/laboutik/,
; ws/printer/, ws/rfid/ (controlvanne), ws/chat/. nginx routes /ws/ and
; /wss/ URLs to this port (see nginx_prod/lespass_prod.conf).
[program:daphne]
command=/home/tibillet/.local/bin/poetry run daphne -b 0.0.0.0 -p 7999 --access-log /DjangoFiles/logs/daphne.logs TiBillet.asgi:application
directory=/DjangoFiles
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/DjangoFiles/logs/daphne.logs
