# supervisor/conf.d/gunicorn.conf
; Gunicorn — serveur HTTP WSGI (port 8002). Même commande que l'ancien
; start.sh, SANS --reload (inutile et coûteux en production).
; nginx route tout le trafic HTTP classique ici.
; / Gunicorn — WSGI HTTP server (port 8002). Same command as the old
; start.sh, WITHOUT --reload. nginx routes all regular HTTP traffic here.
;
; FORMAT DU LOG D'ACCES — trois informations que le format par defaut n'a pas :
;
;   %({x-forwarded-for}i)s  l'IP du VISITEUR. Le defaut (%(h)s) journalise
;                           l'IP de nginx (172.20.0.x) pour 100 % des lignes,
;                           puisque tout le trafic arrive par le proxy.
;   %({host}i)s             le TENANT. Indispensable en multi-tenant : la ligne
;                           de requete ne porte que le chemin, donc deux lieux
;                           differents produisent des lignes identiques.
;   %(M)sms                 la duree de la requete, en millisecondes.
;
; On garde %(r)s (la ligne de requete brute) plutot que de recomposer avec
; %(U)s%(q)s : %(q)s rend la query string SANS son « ? », ce qui donne
; « /event/tag=concert » — lisible comme un chemin, alors que c'en est pas un.
;
; /!\ LES « %% » SONT OBLIGATOIRES. supervisord fait une interpolation Python
; sur la valeur de `command=` : un « %(t)s » simple y est lu comme une variable
; supervisor, elle n'existe pas, et le demarrage echoue sur
; « KeyError: 't' » — gunicorn ne demarre alors PAS DU TOUT. Le doublement
; « %%( » passe un « %( » litteral a gunicorn. Les quotes simples autour du
; format sont tout aussi obligatoires : sans elles, shlex decoupe le format
; sur les espaces et gunicorn ne recoit que le premier morceau.
;
; / ACCESS LOG FORMAT — adds the three things the default lacks: the real
; visitor IP (the default logs nginx's IP for every line), the tenant host
; (without it two tenants log identical lines), and the request duration.
; We keep %(r)s rather than %(U)s%(q)s because %(q)s drops the "?".
; BEWARE: the "%%" are mandatory — supervisord interpolates `command=`, so a
; bare "%(t)s" raises KeyError and gunicorn never starts. Single quotes are
; mandatory too, else shlex splits the format on spaces.
[program:gunicorn]
command=/home/tibillet/.local/bin/poetry run gunicorn TiBillet.wsgi --log-level=info --access-logfile /DjangoFiles/logs/gunicorn.logs --access-logformat '%%({x-forwarded-for}i)s %%({host}i)s %%(t)s "%%(r)s" %%(s)s %%(b)s %%(M)sms "%%(f)s" "%%(a)s"' --log-file /DjangoFiles/logs/gunicorn.logs --error-logfile /DjangoFiles/logs/gunicorn.logs --capture-output -w 18 -b 0.0.0.0:8002
directory=/DjangoFiles
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/DjangoFiles/logs/gunicorn.logs
