# supervisor/conf.d/gunicorn.conf
; Gunicorn — serveur HTTP WSGI (port 8002). Même commande que l'ancien
; start.sh, SANS --reload (inutile et coûteux en production).
; nginx route tout le trafic HTTP classique ici.
; / Gunicorn — WSGI HTTP server (port 8002). Same command as the old
; start.sh, WITHOUT --reload. nginx routes all regular HTTP traffic here.
[program:gunicorn]
command=/home/tibillet/.local/bin/poetry run gunicorn TiBillet.wsgi --log-level=info --access-logfile /DjangoFiles/logs/gunicorn.logs --log-file /DjangoFiles/logs/gunicorn.logs --error-logfile /DjangoFiles/logs/gunicorn.logs --capture-output -w 18 -b 0.0.0.0:8002
directory=/DjangoFiles
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/DjangoFiles/logs/gunicorn.logs
